# packages/chirpstack/luci-app-chirpstack-concentratord/root/etc/uci-defaults/99_set_sx1302_fw_upgrade_command
rule_name=$(uci add luci command)
uci batch << EOI
set luci.$rule_name.name='Upgrade SX1302/3 concentrator firmware'
set luci.$rule_name.command='/opt/libloragw-sx1302/upgrade_fw.sh'
EOI
uci commit
