# packages/chirpstack/luci-app-chirpstack-concentratord/root/etc/uci-defaults/99_set_2g4_fw_upgrade_command
rule_name=$(uci add luci command)
uci batch << EOI
set luci.$rule_name.name='Upgrade 2g4 concentrator firmware'
set luci.$rule_name.command='/opt/libloragw-2g4/upgrade_fw.sh'
EOI
uci commit
